# pi-monitoring/aeris_raspberrypi.service

[Unit]
Description=Aeris iperf test service

[Service]
ExecStart=/home/pi/Aeris-raspberrypi/test_ping.sh

[Install]
WantedBy=multi-user.target
